# src/etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -c /etc/php.ini
stdout_events_enabled=true
stderr_events_enabled=true
